# linux/cmds/fastnode-installer/rootless-system-data/.config/systemd/user/fastnode-updater.service.template
[Unit]
Description=Fastnode Updater

[Service]
Type=oneshot
StandardOutput=null
ExecStart={{.HomeDir}}/.local/share/fastnode/current/fastnode-update "update" "silent"
